# systemd/pi-agent.service
[Unit]
Description=pi-agent Telegram assistant
After=network-online.target awg-quick@awg0.service
Wants=network-online.target awg-quick@awg0.service

[Service]
Type=simple
User=godo
WorkingDirectory=/home/godo/pi-agent
EnvironmentFile=/home/godo/pi-agent/.env
ExecStart=/usr/bin/node /home/godo/pi-agent/dist/index.js
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
